# docker/run/fs/etc/supervisor/conf.d/a2a.conf
[program:a2a]
command=/opt/venv-a0/bin/python -m fast_a2a.server --host 0.0.0.0 --port 9000 --settings /a0/usr/settings.json
environment=HOME="/root",USER="root"
user=root
directory=/
stopwaitsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=5
killasgroup=true
